# systemd/ghostwire-server.service
[Unit]
Description=GhostWire Server
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/ghostwire-server -c /etc/ghostwire/server.toml
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
